# supervisor.conf
[program:elasticsearch]
command=elasticsearch
autostart=true
autorestart=true
startsecs=15
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
user=elasticsearch